# etc/xdg/autostart/HubScreen.desktop
[Desktop Entry]
Type=Application
Name=HubScreen
Comment=Autostart of HubScreen
NoDisplay=false
Exec=surf -F localhost
NotShowIn=GNOME;KDE;XFCE;